# .config/systemd/user/mbsync.timer
[Unit]
Description=Run mbsync every 5 minutes

[Timer]
OnBootSec=30sec
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
